# file_cleaner.cron
0 3 * * * python /srv/salt/scripts/health/rsq_file_cleaner/cleaner.py >> /var/log/rsq_file_cleaner.log 2>&1
